# services/snapserver/finish
#!/bin/sh
echo 'snapserver exited... restarting'
