# services/snapserver/finish
#!/bin/bash
echo 'snapserver exited... restarting'
