# services/snapserver/run
#!/usr/bin/with-contenv /bin/sh
echo 'Starting snapserver'
echo '[stream]' > /root/.config/snapserver/snapserver.conf
# for var in ${!stream*}
# do
#   echo "stream = ${!var}" >> /etc/snapserver.conf
# done
env | grep '^stream_' | cut -d= -f2- | while read -r stream; do 
    echo "stream = ${stream}" >> /root/.config/snapserver/snapserver.conf 
done
echo '[http]' >> /root/.config/snapserver/snapserver.conf 
echo 'doc_root=/snapweb-master/dist' >> /root/.config/snapserver/snapserver.conf 
echo "Configuration:"
cat /root/.config/snapserver/snapserver.conf
exec snapserver

