# services/snapserver/run
#!/usr/bin/with-contenv /bin/sh
echo 'Starting snapserver'
echo '[stream]' > /etc/snapserver.conf
# for var in ${!stream*}
# do
#   echo "stream = ${!var}" >> /etc/snapserver.conf
# done
env | grep '^stream_' | cut -d= -f2- | while read -r stream; do 
	echo "$stream";
    echo "stream = ${stream}" >> /etc/snapserver.conf
done
echo '[http]' >> /etc/snapserver.conf
echo 'doc_root=/snapweb-master/dist' >> /etc/snapserver.conf
echo "Configuration:"
cat /etc/snapserver.conf
exec snapserver

