# services/snapserver/run
#!/usr/bin/with-contenv /bin/bash
echo 'Starting snapserver'
echo '[stream]' > /etc/snapserver.conf
for var in ${!stream*}
do
  echo "source = ${!var}" >> /etc/snapserver.conf
done
echo "Configuration:"
cat /etc/snapserver.conf
exec snapserver
